# root/etc/s6/search/finish
#!/usr/bin/env sh
exec /bin/s6-svscanctl -t /etc/s6
